# debian/packagekit.postrm
#!/bin/sh
# postrm script for packagekit

set -e

#DEBHELPER#

if [ "$1" = "purge" ]; then

        if [ -x /etc/init.d/dbus ]; then
                invoke-rc.d dbus force-reload || true
        fi

fi

